# files/etc/init/mount-extra.conf
# mount-extra - Mount extra filesystems after creating required directories
#

description	"Mount local extra filesystems"

start on mounted

task

script
   /usr/local/sbin/mount-extra
end script
